# default-settings/usr/share/distro-transmute/default-settings/base-files/etc/profile.d/soulmate-path.sh
for p in /usr/local/sbin /usr/sbin /sbin ; do
	if ! echo ":$PATH:" |grep -q ":$p:" ; then
		export PATH="$p:$PATH"
	fi
done

