# tg3-fix/etc/udev/rules.d/80-tg3-fix.rules
ACTION=="add", SUBSYSTEM=="net", ATTRS{vendor}=="0x14e4", ATTRS{device}=="0x1687", RUN+="/usr/share/tg3-fix/tg3-udev-fix %k"
